# 05-ansible-bootstrap/roles/juicefs/templates/juicefs-mount.service.j2
[Unit]
Description=JuiceFS Mount Service (MDI Managed)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# 환경 변수 설정
Environment="ACCESS_KEY={{ minio_root_user }}"
Environment="SECRET_KEY={{ minio_root_password }}"

# 마운트 포인트가 존재하지 않으면 자동으로 생성하도록 설정합니다.
ExecStartPre=/usr/bin/mkdir -p {{ storage_info.mount_point }}

# 메타데이터 주소만으로도 MinIO 백엔드를 자동으로 찾아 마운트합니다.
ExecStart={{ juicefs_bin_path }} mount \
    --cache-size {{ storage_info.cache_size }} \
    -o allow_other \
    "redis://:{{ redis_password }}@{{ storage_info.redis_url }}/0" \
    "{{ storage_info.mount_point }}"

ExecStop=/usr/bin/umount -l {{ storage_info.mount_point }}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target